# nginx-proxy-manager/scripts/mailcow-cert-update.timer
[Unit]
Description=Update Mailcow SSL certificate

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

